# root/etc/crontabs/root
# do daily/weekly/monthly maintenance
# min   hour    day     month   weekday command
*/15    *       *       *       *       run-parts /etc/periodic/15min
0       *       *       *       *       run-parts /etc/periodic/hourly
0       2       *       *       *       run-parts /etc/periodic/daily
0       3       *       *       6       run-parts /etc/periodic/weekly
0       5       1       *       *       run-parts /etc/periodic/monthly

# scrutiny collector
0 0 * * * /usr/local/bin/scrutiny-collector-metrics run --api-endpoint {SCRUTINY_API_ENDPOINT} >> /config/log/scrutiny-collector-metrics.log 2>&1
