# root/etc/services.d/web/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/local/bin/scrutiny start
